# php-worker/supervisord.d/outgoing-api-worker.conf
[program:outgoing-api-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/outgoing_api/artisan queue:work --tries=3
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/outgoing-api-worker.log
stopwaitsecs=3600
